# config/mpd/dot-config/systemd/user/mpd-updatedb.service
[Unit]
Description=MPD Trigger Database Update
Requires=mpd.service

[Service]
Type=oneshot
ExecStart=mpc update
