# config/mpd/dot-config/systemd/user/mpd-mk-socket-dir.service
[Unit]
Description=MPD Create Socket Directory

[Service]
Type=oneshot
ExecStart=mkdir -p /run/user/%U/mpd
