# systemd/.config/systemd/user/mbsync.service
[Unit]
Description=Mail sync
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=mbsync -a
