# config/borg/dot-config/systemd/user/borg.timer
[Unit]
Description=Hourly Borg Backup

[Timer]
OnCalendar=hourly
AccuracySec=1minute
Presistent=true

[Install]
WantedBy=timers.target
