# config/borg/dot-config/systemd/user/borg-prune.timer
[Unit]
Description=Daily Borg Backup Pruning

[Timer]
OnCalendar=daily
AccuracySec=30minute
Persistent=true

[Install]
WantedBy=timers.target
